# dist/hyperic-hq-installer/rpm-automation/classic/RPMs/support_files_EE/hyperic-hq-agent.init.suse
#! /bin/sh
#
# chkconfig: 2345 80 30
# description: Starts and stops the hyperic agent.
#
# /etc/init.d/hyperic-hq-agent
#
### BEGIN INIT INFO
# Provides:          hyperic-hq-agent
# Required-Start:    $network
# X-UnitedLinux-Should-Start: $time 
# Required-Stop:     
# X-UnitedLinux-Should-Stop: 
# Default-Start:     3 5
# Default-Stop:      0 1 2 6
# Short-Description: hyperic-hq-agent init
# Description:       Starts Hyperic HQ Agent, located by default in
#		/opt/hyperic/hyperic-hq-agent directory.
### END INIT INFO


# Check for missing binaries (stale symlinks should not happen)
AGENT_DIR=/opt/hyperic/hyperic-hq-agent
test -x $AGENT_DIR || exit 5

# Set default start-up user
USER=hyperic

# Set HQ_JAVA_HOME if you want to use a different JDK
export HQ_JAVA_HOME=/usr/java/jdk1.5.0_12

# Shell functions sourced from /etc/rc.status:
#      rc_check         check and set local and overall rc status
#      rc_status        check and set local and overall rc status
#      rc_status -v     ditto but be verbose in local rc status
#      rc_status -v -r  ditto and clear the local rc status
#      rc_status -s     display "skipped" and exit with status 3
#      rc_status -u     display "unused" and exit with status 3
#      rc_failed        set local and overall rc status to failed
#      rc_failed <num>  set local and overall rc status to <num>
#      rc_reset         clear local rc status (overall remains)
#      rc_exit          exit appropriate to overall rc status
#      rc_active	checks whether a service is activated by symlinks
#      rc_splash arg    sets the boot splash screen to arg (if active)
. /etc/rc.status

# Reset status of this service
rc_reset

# Return values acc. to LSB for all commands but status:
# 0	  - success
# 1       - generic or unspecified error
# 2       - invalid or excess argument(s)
# 3       - unimplemented feature (e.g. "reload")
# 4       - user had insufficient privileges
# 5       - program is not installed
# 6       - program is not configured
# 7       - program is not running
# 8--199  - reserved (8--99 LSB, 100--149 distrib, 150--199 appl)
# 
# Note that starting an already running service, stopping
# or restarting a not-running service as well as the restart
# with force-reload (in case signaling is not supported) are
# considered a success.

case "$1" in
    start)
	echo -n "Starting hyperic agent "
	/bin/su $USER -c "$AGENT_DIR/hq-agent.sh start > /dev/null 2>&1"

	# Remember status and be verbose
	rc_status -v
	;;
    stop)
	echo -n "Shutting down hyperic agent "
        /bin/su $USER -c "$AGENT_DIR/hq-agent.sh stop > /dev/null 2>&1"

	# Remember status and be verbose
	rc_status -v
	;;
    restart)
	## Stop the service and regardless of whether it was
	## running or not, start it again.
	$0 stop
	$0 start

	# Remember status and be quiet
	rc_status
	;;
    *)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
	;;
esac
rc_exit
